# root/etc/services.d/hera/run
#!/usr/bin/execlineb

foreground {
  hera
}
